# scripts/com.mealprepper.plan-week.plist.example
<!-- macOS launchd example only — on Linux use scripts/install_systemd.sh -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mealprepper.plan-week</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-m</string>
    <string>mealprepper</string>
    <string>plan-week</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/MealPrepper</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>6</integer>
    <key>Hour</key>
    <integer>10</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
</dict>
</plist>
